# etc/supervisor/conf.d/nscd.conf

[program:nscd]
command=/usr/sbin/nscd -F
autostart=true
autorestart=true

