# etc/supervisor/conf.d/rsyslog.conf

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true

